# go/wgetbash.service.example
[Unit]
Description=wgetbash
After=network.target

[Service]
Type=simple
User=root
ExecStart=/opt/wgetbash/wgetbash
Restart=always
RestartSec=5

Environment=PORT=8000
Environment=DB_PATH=/data/wgetbash.db
Environment=AUTH_URL=https://auth-center.example.com
Environment=AUTH_INTERNAL=http://localhost:8886
Environment=APP_URL=https://wgetbash.example.com
Environment=APP_TOKEN=wgetbash:your-token-here
Environment=SECRET_KEY=your-secret-key-here

[Install]
WantedBy=multi-user.target
